# tg-forwarder.service
[Unit]
Description=Telegram Message Forwarder Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/user/tg-forwarder/forwarder.py /home/user/tg-forwarder/config.json
WorkingDirectory=/home/user/tg-forwarder
Restart=always
RestartSec=10
StandardOutput=append:/home/user/tg-forwarder/forwarder.log
StandardError=append:/home/user/tg-forwarder/forwarder.log

[Install]
WantedBy=default.target
